# systemd/user/20-20-20.service
[Unit]
Description=20-20-20 Rule

[Service]
ExecStart=/bin/sh /home/ConnorM/scripts/20-20-20
